# rootfs/base/etc/systemd/system/snapraid-sync.service
[Unit]
Description=SnapRAID Sync

[Service]
Type=oneshot
User=root
ExecStart=snapraid sync
