# rootfs/base/etc/systemd/system/snapraid-scrub.service
[Unit]
Description=SnapRAID Scrub

[Service]
Type=oneshot
User=root
ExecStart=snapraid scrub
