# rootfs/base/etc/systemd/system/snapraid-sync.timer
[Unit]
Description=Runs SnapRAID sync nightly

[Timer]
OnCalendar=Mon-Sun, 3:00

[Install]
WantedBy=timers.target
