# rootfs/base/etc/systemd/system/snapraid-scrub.timer
[Unit]
Description=Runs SnapRAID scrub weekly

[Timer]
OnCalendar=Mon, 5:00

[Install]
WantedBy=timers.target
